# bettercap.service.template
[Unit]
Description=bettercap WiFi Sentinel (production - USB wlan1mon, API forced early)
After=network-online.target NetworkManager.service
Wants=network-online.target
[Service]
Type=simple
User=root
WorkingDirectory=/home/joel/wifi-sentinel
ExecStartPre=/home/joel/wifi-sentinel/fix-pi-bettercap.sh --no-start
ExecStart=/usr/bin/bettercap -iface wlan1mon -eval "set api.rest.address 0.0.0.0; set api.rest.port 8081; api.rest on" -caplet /home/joel/wifi-sentinel/wifi-sentinel.cap --no-colors
Restart=always
RestartSec=8
TimeoutStartSec=60
StandardOutput=journal
StandardError=journal
[Install]
WantedBy=multi-user.target
